# files/etc/profile.d/cassandra.sh
export PATH=$PATH:/opt/apache-cassandra/bin
